# etc/supervisor-ansible.conf
[program:ansible-container]
stopwaitsecs=30
directory=/home/serviceapi/app
command=/usr/local/bin/ansible-container run --production
